# dist/systemd/aws-cloudwatch-daemon.service
[Unit]
Description=Push instance metrics to CloudFront

[Service]
ExecStart=/usr/bin/docker run --rm -v /usr/share/ca-certificates:/etc/ssl/certs:ro yieldr/aws-cloudwatch-daemon

[Install]
WantedBy=multi-user.target
